# scripts/com.gregehrenberg.thebat-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gregehrenberg.thebat-tracker</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/gregehrenberg/MLB/thebat_tracker/scripts/run_daily.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/gregehrenberg/MLB/thebat_tracker/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/gregehrenberg/MLB/thebat_tracker/logs/launchd.err.log</string>
</dict>
</plist>
